# roles/terminal/files/launchagents/com.stefanroex.bm-cmux-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stefanroex.bm-cmux-status</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/stefanroex/.local/bin/bm-cmux-status</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- bb lives in /opt/homebrew/bin, lsof (used by REPL discovery) in /usr/sbin -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartInterval</key>
  <integer>10</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/bm-cmux-status.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/bm-cmux-status.log</string>
</dict>
</plist>
